# 5.3/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/start-apache2
numprocs=1
autostart=true
autorestart=true

[program:mysqld]
command=/start-mysqld
numprocs=1
autostart=true
autorestart=true
